# files/auto-shutdown.service
[Unit]
Description=Auto-shutdown service for idle AWS instances
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/auto-shutdown.sh start
ExecStop=/usr/local/bin/auto-shutdown.sh stop
Restart=always
RestartSec=30
Environment=IDLE_TIME_MINUTES=60
Environment=CHECK_INTERVAL_MINUTES=5

[Install]
WantedBy=multi-user.target 